# docker/frankenphp/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid

[program:laravel-frankenphp]
command=php artisan octane:start --server=frankenphp --host=0.0.0.0 --port=80 --admin-port=2019
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/frankenphp.log
stderr_logfile=/var/log/supervisor/frankenphp-error.log
stdout_logfile_maxbytes=100MB
stderr_logfile_maxbytes=100MB

[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work redis --tries=5 --timeout=180
autostart=true
autorestart=true
numprocs=3
stdout_logfile=/var/log/supervisor/laravel-queue.log
stderr_logfile=/var/log/supervisor/laravel-queue-error.log
stdout_logfile_maxbytes=100MB
stderr_logfile_maxbytes=100MB

[program:laravel-schedule]
command=/bin/sh -c "while true; do php /var/www/artisan schedule:run --no-interaction >> /var/log/supervisor/schedule.log 2>&1; sleep 60; done"
autostart=true
autorestart=true
startsecs=0
numprocs=1
stdout_logfile=/var/log/supervisor/laravel-schedule.log
stderr_logfile=/var/log/supervisor/laravel-schedule-error.log
stdout_logfile_maxbytes=100MB
stderr_logfile_maxbytes=100MB

[program:laravel-reverb]
command=php artisan reverb:start --host=0.0.0.0 --port=6001
autostart=true
autorestart=true
priority=10
startsecs=15
stdout_logfile=/var/log/supervisor/laravel-reverb.log
stderr_logfile=/var/log/supervisor/laravel-reverb-error.log
stdout_logfile_maxbytes=100MB
stderr_logfile_maxbytes=100MB
